# exercises/services-templates-and-registering-variables/solution/node-app.conf
description "daemon for node app"

start on startup
stop on shutdown
respawn

script
  {{ which_nodejs.stdout | default('/usr/bin/nodejs') }} {{ node_app_path }}
end script
