# exercises/08-managing-secrets/app.conf
description "daemon for node app"

start on startup
stop on shutdown
respawn

env API_KEY={{api_key}}

script
  {{node_path.stdout}} /mnt/app.js
end script
